# software/bedrock-connect/assets/services.ini
[supervisord]
nodaemon=true

[program:bedrock-connect]
directory=/home/container/bedrock-connect
command=/bin/bash -c "java -Xms128M -Xmx500M -jar BedrockConnect.jar"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:mcxboxbroadcast]
directory=/home/container/mcxboxbroadcast
command=/bin/bash -c "java -Xms128M -Xmx500M -jar MCXboxBroadcastStandalone.jar"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0