# software/bedrock-connect/assets/services.ini
[supervisord]
nodaemon=false

[program:bedrock-connect]
directory=/home/container/bedrock-connect
command=/bin/bash -c "java -Xms128M -Xmx${SERVER_MEMORY}M -jar BedrockConnect.jar port=${SERVER_PORT} custom_servers=/home/container/bedrock-connect/custom_servers.json featured_server_ips=/home/container/bedrock-connect/featured_server_ips.json language=/home/container/bedrock-connect/language.json"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:mcxboxbroadcast]
directory=/home/container/mcxboxbroadcast
command=/bin/bash -c "java -Xms128M -Xmx${SERVER_MEMORY}M -jar MCXboxBroadcast.jar"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0