# software/bedrock-connect/assets/services.ini
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface


[program:bedrock-connect]
directory=/home/container/bedrock-connect
command=/bin/bash -c "java -Xms128M -Xmx${SERVER_MEMORY}M -jar BedrockConnect.jar port=${SERVER_PORT}"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:mcxboxbroadcast]
directory=/home/container/mcxboxbroadcast
command=/bin/bash -c "java -Xms128M -Xmx${SERVER_MEMORY}M -jar MCXboxBroadcast.jar"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0